# systemd/hyrovi-touch-admin.service
[Unit]
Description=Hyrovi Touch Panel Admin
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=hyrovi-panel
Group=hyrovi-panel
WorkingDirectory=/opt/hyrovi-touch-panel/current
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/hyrovi-touch-panel/venv/bin/python /opt/hyrovi-touch-panel/current/admin/server.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
